# rootfs/etc/cron.d/osticket-cron
* * * * * root /opt/bitnami/php/bin/php /app/api/cron.php > /var/log/apicron 2>&1
